# user/services/shell.service
[Unit]
Description=Shell

[Service]
Type=simple
ExecStart=/bin/NovaShell
Restart=always
ExecStartPre=-/bin/about.elf
